# web/ansible/roles/webservers/templates/etc/supervisor/conf.d/mydota.conf
; Managed by Ansible. Do not edit manually.

[program:mydota]
environment=MYDOTA_CONFIG="{{ appcontainers_gopath  }}/src/hamak/mydota/config/live.json"
directory=/var/apps/mydota/gopath/bin
command=/var/apps/mydota/gopath/bin/mydota
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mydota.out.log
stdout_logfile_maxbytes=50MB
